# systemd/bible-clock.service
[Unit]
Description=Bible Clock E-ink Display Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/Bible-Clock-v3
Environment=PATH=/home/admin/Bible-Clock-v3/venv/bin
ExecStart=/home/admin/Bible-Clock-v3/venv/bin/python3 /home/admin/Bible-Clock-v3/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/admin/Bible-Clock-v3 /var/log

[Install]
WantedBy=multi-user.target